# config/systemd/mcp-dev-platform.service
[Unit]
Description=MCP-Server Development Platform
Documentation=https://github.com/your-org/mcp-dev-platform
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/mcp-dev-platform
Environment=PYTHONPATH=/home/pi/mcp-dev-platform
Environment=PORT=8000
Environment=HOST=0.0.0.0
ExecStart=/usr/bin/python3 main.py --host 0.0.0.0 --port 8000
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mcp-dev-platform

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/home/pi/mcp-dev-platform/projects /home/pi/mcp-dev-platform/logs

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target 